# install/pass-secrets.service.in
[Unit]
Description=Pass SecretService

[Service]
Type=dbus
BusName=org.freedesktop.secrets
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/pass-secrets

[Install]
WantedBy=default.target
Alias=dbus-org.freedesktop.secrets.service